# crontab
# m h  dom mon dow   command
45 3 * * 6 /bin/bash ~/WORK/test-build2/all_start.sh
